# packaging/sentinel-backend.service
[Unit]
Description=Sentinel Biometric Backend Service
After=network.target video.target

[Service]
Type=simple
User=root
Group=root
# Camera access: add video group so /dev/video* devices are accessible
SupplementaryGroups=video

# Environment for TensorFlow/Mesa silencing
Environment=TF_CPP_MIN_LOG_LEVEL=3
Environment=GLOG_minloglevel=3
Environment=LIBGL_DEBUG=quiet
Environment=MESA_DEBUG=silent

# Create /run/sentinel directory automatically
RuntimeDirectory=sentinel
RuntimeDirectoryMode=0755

# Log directory (sentinel writes logs here — root:root 750)
LogsDirectory=sentinel
LogsDirectoryMode=0750

WorkingDirectory=/usr/lib/project-sentinel
ExecStart=/usr/lib/project-sentinel/venv/bin/sentinel-daemon

# Restart logic
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
